# deploy/fnup-api.service
[Unit]
Description=FoldersNUsersParser FastAPI backend
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/FoldersNUsersParser
EnvironmentFile=/srv/FoldersNUsersParser/.env
ExecStart=/srv/FoldersNUsersParser/.venv/bin/python -m uvicorn backend.app.main:app --host 127.0.0.1 --port 8002 --workers 1
Restart=always
RestartSec=10
User=softportal
Group=softportal

[Install]
WantedBy=multi-user.target
